# service_template.service
; linux_service_template.service
; Template unit: run as an instance service like:
;   sudo systemctl start robot-app@robotop
; IMPORTANT: save this as /etc/systemd/system/robot-app@.service

[Unit]
Description=Robot App (Example) for user %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i
ExecStart=/home/%i/robotics/scripts/robot_app.sh
Restart=always
RestartSec=2

; Send stdout/stderr to the journal so journalctl shows your script output
StandardOutput=journal
StandardError=journal

; Optional: give the service a clean, predictable environment
; Environment=ROBOT_MODE=production
; Environment=ROBOTICS_WS=/home/%i/robotics

[Install]
WantedBy=multi-user.target

